# debian/overlay/etc/systemd/system/serial-getty@ttyPS0.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear %I --keep-baud 115200 linux
